# MY_REPOS/misc-experiments/_FIREBFIRE/chef-google_auth_proxy/templates/default/upstart.conf.erb
description     "Google Auth Proxy instance"

start on runlevel [2345]
stop on shutdown

respawn

console log

setuid <%= @user %>
exec <%= node['google_auth_proxy']['bin_path'] %>/google_auth_proxy -config=/etc/google_auth_proxy/<%= @service_name %>.conf
